# config/ebay-bot.service
[Unit]
Description=eBay Discord Bot Manager
After=network.target

[Service]
Type=simple
User=jacob
WorkingDirectory=/home/jacob/ebay-bot
ExecStart=/home/jacob/ebay-bot/ebaymanager-bot-linux
Restart=always
RestartSec=10
StandardOutput=append:/home/jacob/ebay-bot/bot.log
StandardError=append:/home/jacob/ebay-bot/bot-error.log

# Environment settings
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target
